# stock-pump-railway-legacy.service
[Unit]
Description=Railway Legacy Stock Pump Simulation Bot v32.1
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/fre-stock-alarm-bot-railway
EnvironmentFile=/home/ubuntu/fre-stock-alarm-bot-railway/.env
ExecStart=/home/ubuntu/fre-stock-alarm-bot-railway/.venv/bin/python /home/ubuntu/fre-stock-alarm-bot-railway/stock_pump_railway_legacy_v32_1.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
